# docker-rmeditor/rmeditor.supervisord.conf
[program:rmeditor]
environment=PORT=8000,REDISTOGO_URL=%(ENV_REDIS_PORT_6379_TCP)s
directory=/opt/rmeditor/current
command=/opt/nodejs/current/bin/node server.js
